# init/throttle-cpu-temp.service
[Unit]
Description=throttle cpu temp

[Service]
ExecStart=/usr/local/bin/throttle-cpu-temp 70

[Install]
WantedBy=multi-user.target
